# templates/program.conf.j2
[program:{{ item.name }}]
command={{ item.command }}
{% if item.configuration is defined %}
{% for k, v in item.configuration.iteritems() %}
{{k}}={{v}}
{% endfor %}
{% endif %}
